# init/playlists.conf
start on started networking
script
  su -l barry /opt/sb2m3u/convert_playlists.sh
end script
